# com.codex.ndx-shadow-1310.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.ndx-shadow-1310</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>/Users/cuichi/Library/Application Support/fund_tracker_launch/run_ndx_shadow_1310.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>13</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>TZ</key>
    <string>Asia/Singapore</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/cuichi/Library/Logs/fund_tracker/ndx-shadow-1310.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/cuichi/Library/Logs/fund_tracker/ndx-shadow-1310-error.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
